# examples/rank-tracker.timer
[Unit]
Description=Daily rank tracker check

[Timer]
OnCalendar=*-*-* 06:00:00
Persistent=true

[Install]
WantedBy=timers.target
